# meta/config/com.juliaz.start-system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.juliaz.start-system</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/raphael/juliaz_agents/start-system.sh</string>
    </array>

    <!-- Run once at login -->
    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/raphael/juliaz_agents</string>

    <key>StandardOutPath</key>
    <string>/Users/raphael/juliaz_agents/logs/launchagent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/raphael/juliaz_agents/logs/launchagent-error.log</string>

    <!-- Don't restart — this is a one-shot boot script -->
    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/raphael</string>
    </dict>
</dict>
</plist>
